# install/root/root/lib/crontab
# min	hour	dom	month	dow	command
39	1,13	*	*	*	/home/ec2-user/lib/gww/script/suppress_output_on_success certbot renew --no-self-upgrade
