# install/root/root/lib/crontab
# min	hour	dom	month	dow	command
39	1,13	*	*	*	certbot renew --no-self-upgrade
